# dist/etc/systemd/system/goordinator.service
[Unit]
Description=Goordinator
Wants=network-online.target
After=network-online.target

[Service]
ExecStart=/usr/local/bin/goordinator --cfg-file /etc/goordinator/config.toml
ProtectSystem=strict

[Install]
WantedBy=multi-user.target
